# src/ci/build/supervisord.conf
[program:portal]
command=npm run start
directory=/app/portal/ai-verify-portal
autostart=true
autorestart=true
stderr_logfile=/var/log/portal.err.log
stdout_logfile=/var/log/portal.out.log
user=root

[program:test-engine-app]
command=python3 -m test_engine_app
directory=/app/test-engine/app
autostart=true
autorestart=true
stderr_logfile=/var/log/test-engine-app.err.log
stdout_logfile=/var/log/test-engine-app.out.log
user=root
